# contrib/clipd.service
[Unit]
Description=clipto clipboard daemon
# Start after the Wayland session is available so WAYLAND_DISPLAY is set.
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/clipd
Restart=on-failure
RestartSec=2

# Seal the encryption key to this machine via the TPM. Generate with:
#   dd if=/dev/urandom bs=32 count=1 2>/dev/null \
#     | systemd-creds encrypt --name=clipto-key - \
#       %h/.config/clipto/clipto-key.cred
LoadCredentialEncrypted=clipto-key:%h/.config/clipto/clipto-key.cred

[Install]
WantedBy=graphical-session.target
